# mariadb/mariadb-10.4/supervisor/fluentd.conf
[program:fluentd]
command			= td-agent
autostart		= true
autorestart		= true
redirect_stderr	= true 
